# debian/rpi/postinst
#!/bin/sh
set -e

. /usr/share/debconf/confmodule

if [ "$1" = "configure" ]; then
    db_get rws2811/enable-spi
    if [ "$RET" = "true" ]; then
        # Find the boot config file
        if [ -f /boot/firmware/config.txt ]; then
            CONFIG=/boot/firmware/config.txt
        elif [ -f /boot/config.txt ]; then
            CONFIG=/boot/config.txt
        else
            echo "Warning: could not find Raspberry Pi boot config" >&2
            exit 0
        fi

        if grep -q '^dtparam=spi=on' "$CONFIG"; then
            echo "SPI is already enabled in $CONFIG"
        elif grep -q '^#dtparam=spi=on' "$CONFIG"; then
            sed -i 's/^#dtparam=spi=on/dtparam=spi=on/' "$CONFIG"
            echo "SPI enabled in $CONFIG (uncommented). Reboot required."
            touch /var/run/reboot-required
            echo "rws2811: SPI was enabled in boot config" >> /var/run/reboot-required.pkgs
        else
            echo "dtparam=spi=on" >> "$CONFIG"
            echo "SPI enabled in $CONFIG (added). Reboot required."
            touch /var/run/reboot-required
            echo "rws2811: SPI was enabled in boot config" >> /var/run/reboot-required.pkgs
        fi
    fi

    db_get rws2811/buffer-size
    BUFFER_SIZE="$RET"
    db_get rws2811/spi-speed
    SPI_SPEED="$RET"
    db_get rws2811/spi-device
    SPI_DEVICE="$RET"
    db_get rws2811/frame-rate
    FRAME_RATE="$RET"
    db_get rws2811/listen-address
    LISTEN_ADDR="$RET"

    mkdir -p /etc/default
    cat > /etc/default/rws2811 <<EOF
# Configuration for rws2811 - generated by package install
# Edit this file to change runtime options, then restart the service.
# The SPI device is set via the systemd instance name, e.g.:
#   systemctl enable --now rws2811@spidev0.0
RWS2811_OPTS="--size $BUFFER_SIZE --speed $SPI_SPEED --rate $FRAME_RATE --listen $LISTEN_ADDR"
EOF

    # Extract the device node name (e.g. spidev0.0) from the full path
    SPI_INSTANCE="${SPI_DEVICE#/dev/}"
    if [ ! -e "$SPI_DEVICE" ]; then
        echo "Warning: $SPI_DEVICE does not exist. The service will not be enabled." >&2
        echo "Once the device is available, enable it manually with:" >&2
        echo "  systemctl enable --now rws2811@${SPI_INSTANCE}" >&2
    elif command -v systemctl >/dev/null 2>&1; then
        systemctl enable --now "rws2811@${SPI_INSTANCE}"
    fi
fi

#DEBHELPER#
